# deployment/harvester.timer
[Unit]
Description=umwelt.info harvester timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
