# deployment/harvester-daily.timer
[Unit]
Description=umwelt.info harvester timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=7200
Persistent=true

[Install]
WantedBy=default.target
