# deployment/harvester-weekly.timer
[Unit]
Description=umwelt.info harvester timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=7200
Persistent=true

[Install]
WantedBy=default.target
